# folders/etc/systemd/system/greetd.service.d/wait-for-flatpaks.conf
# First-boot only: ensures greetd does not show the login screen until
# Flatpaks prefetched during installation have been deployed. On subsequent
# boots, ditana-flatpak-finalize.service is skipped (ConditionPathExists
# fails when /var/lib/ditana/flatpak-pending is absent) and this drop-in
# becomes a no-op.
[Unit]
After=ditana-flatpak-finalize.service
