# folders/etc/systemd/system/greetd.service.d/wait-for-flatpaks.conf
[Unit]
After=ditana-flatpak-finalize.service
Wants=ditana-flatpak-finalize.service
