# folders/usr/lib/systemd/user/ditana-wayfire-xdg-autostart.service
[Unit]
# =======================================================================================
# == Ditana | XDG Autostart Bridge for Wayfire Sessions                              ==
# =======================================================================================
#
# This service is the second part of our bridge between Wayfire and the
# systemd-based XDG autostart mechanism.
#
# PURPOSE:
# Its sole purpose is to act as a valid "anchor" for the systemd target
# `xdg-desktop-autostart.target`.
#
# BACKGROUND:
# For security reasons, the `xdg-desktop-autostart.target` is configured to refuse
# being started manually. It can only be pulled in as a dependency by another unit.
#
# HOW IT WORKS:
# The `Wants=` directive declares a non-critical dependency. As soon as this service
# (`ditana-wayfire-xdg-autostart.service`) is started, systemd sees this directive
# and subsequently starts `xdg-desktop-autostart.target` as well. The `Before=`
# directive helps ensure a clean startup order.
#
Description=Ditana | XDG Autostart Bridge for Wayfire Sessions
Wants=xdg-desktop-autostart.target
Before=xdg-desktop-autostart.target

[Service]
# The actual work – setting the environment variables – is handled by the
# `dbus-update-activation-environment` command in `wayfire.ini` before this service
# is even called. Therefore, this service itself needs to do nothing more than
# start successfully and exit immediately.
Type=oneshot
ExecStart=/bin/true