# folders/etc/systemd/system/ditana-flatpak-finalize.service
[Unit]
Description=Finalize prefetched Flatpak applications (Ditana first-boot)
Documentation=https://github.com/acrion/ditana-installer
After=network-online.target
Wants=network-online.target
ConditionPathExists=/var/lib/ditana/flatpak-pending
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/lib/ditana/flatpak-finalize
RemainAfterExit=no
StandardOutput=journal
StandardError=journal
TimeoutStartSec=15min

[Install]
WantedBy=multi-user.target
