# folders/etc/systemd/system/ditana-initialize-system.service
[Unit]
Description=Post install steps for Ditana
Wants=network-online.target
After=network-online.target NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/lib/ditana/initialize-system-as-root
ExecStartPost=/bin/systemctl disable ditana-initialize-system.service

[Install]
WantedBy=multi-user.target
