# folders/etc/systemd/system/ditana-rate-mirrors.service
[Unit]
Description=Update pacman mirrorlist via rate-mirrors

[Service]
Type=oneshot
ExecStart=/usr/bin/rate-mirrors --allow-root --save=/etc/pacman.d/mirrorlist arch
