# folders/etc/systemd/system/ditana-flatpak-finalize.service
[Unit]
Description=Finalize prefetched Flatpak applications (Ditana first-boot)
Documentation=https://github.com/acrion/ditana-installer
After=network-online.target
Wants=network-online.target
ConditionPathExists=/var/lib/ditana/flatpak-pending
Before=greetd.service

[Service]
Type=oneshot
ExecStartPre=/usr/lib/ditana/flatpak-show-progress start
ExecStart=/usr/lib/ditana/flatpak-finalize
ExecStartPost=/usr/lib/ditana/flatpak-show-progress stop
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
