# folders/etc/systemd/system/ditana-rate-mirrors.timer
[Unit]
Description=Daily mirrorlist update

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target