# folders/etc/systemd/user/foot-server.service.d/ditana-shutdown.conf
# foot-server holds onto Wayland connection during shutdown until SIGKILL.
# When the compositor dies (logout), foot-server's cleanup hangs trying
# to roundtrip a dead Wayland socket, blocking greetd restart for 90s.
# Cap the stop timeout aggressively — by the time SIGTERM arrives, the
# compositor is already gone and clean shutdown is impossible anyway.
[Service]
TimeoutStopSec=3s
