# folders/etc/systemd/system/systemd-zram-setup@zram0.service.d/ditana-mem-limit.conf
# Physical memory cap for the ZRAM device (safety backstop).
#
# ZRAM's disksize only limits the amount of *original* (uncompressed)
# data; the physical footprint is original/ratio and therefore unbounded
# in the pathological case of incompressible memory contents. This cap
# guarantees ZRAM never consumes more than 75% of physical RAM. When the
# cap is reached, further swap writes to the device fail and the kernel
# will eventually OOM-kill the largest consumer -- a recoverable failure
# mode, unlike total memory exhaustion (system freeze).
#
# Note: this is an emergency brake, not a graceful spillover. Graceful
# spillover to a disk swap partition happens only when the ZRAM disksize
# itself is exhausted.
[Service]
ExecStartPost=/bin/sh -c 'echo @@ZRAM_MEM_LIMIT_GIB@@G > /sys/block/%i/mem_limit'
